# Red/custom_service/kerna1.service
[Unit]
Description=Kernal_Things

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/kernal

[Install]
WantedBy=multi-user.target
